# deploy/kam-transfer.service
# kam-transfer companion daemon — systemd unit
#
# Install (e.g. on a Raspberry Pi):
#   sudo cp kam-transfer-mtp-linux-arm64 /usr/local/bin/kam-transfer
#   sudo cp deploy/kam-transfer.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kam-transfer
#
# Adjust User= to taste. On first run the daemon seeds a config at
# ~User/.config/kam-transfer/config.yaml (i.e. /home/pi/... for User=pi);
# edit it for the listen port / auth token, then
# `systemctl restart kam-transfer`. Override the path with --config or
# the KAM_TRANSFER_CONFIG environment variable.
#
# For the optional Display HAT Mini status screen, the service user must
# be in the spi/i2c/gpio groups (`sudo usermod -aG spi,i2c,gpio pi`), or
# uncomment the SupplementaryGroups line below.

[Unit]
Description=KAM Mission Planner companion daemon
Documentation=https://github.com/kamdynamics/kam-transfer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kmendez
WorkingDirectory=/home/kmendez
ExecStart=/usr/local/bin/kam-transfer serve
Restart=on-failure
RestartSec=5
# Mounting camera/drone USB storage (a card reader, or a drone in Mass
# Storage mode) read-only needs CAP_SYS_ADMIN. Without it the daemon
# still runs — it just can't browse media off USB block devices.
AmbientCapabilities=CAP_SYS_ADMIN
# Display HAT Mini status screen — grant GPIO/SPI/I2C access:
#SupplementaryGroups=spi i2c gpio

[Install]
WantedBy=multi-user.target
