# deploy/kam-transfer.service
# kam-transfer companion daemon — systemd unit
#
# Install (e.g. on a Raspberry Pi):
#   sudo cp kam-transfer-mtp-linux-arm64 /usr/local/bin/kam-transfer
#   sudo cp deploy/kam-transfer.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kam-transfer
#
# Adjust User= to taste. On first run the daemon seeds a config at
# ~User/.config/kam-transfer/config.yaml (i.e. /home/pi/... for User=pi);
# edit it for the listen port / auth token, then
# `systemctl restart kam-transfer`. Override the path with --config or
# the KAM_TRANSFER_CONFIG environment variable.

[Unit]
Description=KAM Mission Planner companion daemon
Documentation=https://github.com/kamdynamics/kam-transfer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/local/bin/kam-transfer serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
